# skills/garden-tracker/scripts/dashboard/kuro-dashboard.service
[Unit]
Description=UMI Dashboard (Garden DB)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/ubuntu/.openclaw/workspace/skills/garden-tracker/scripts/dashboard/server.py
WorkingDirectory=/home/ubuntu/.openclaw/workspace/skills/garden-tracker/scripts/dashboard
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
